# roles/vpnmgr-agent/templates/vpn-mgr.supervisor.conf
{% if 0 %}
[program:vpnmgr-agent-web]
directory = /opt/vpnmgr-agent/
environment = VPNMGR_WEB_API_KEY="{{vpn_mgr_api_key}}",VPNMGR_WEB_API_DEBUG="0"
command = /opt/vpnmgr-agent/env/bin/vpnmgr-web 7080
stderr_logfile = /var/log/supervisor/%(program_name)s_error.log
stdout_logfile = /var/log/supervisor/%(program_name)s.log
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
{% endif %}

[program:vpnmgr-agent]
directory = /opt/vpnmgr-agent/
command = /opt/vpnmgr-agent/env/bin/vpnmgr-agent --id "{{ vpnmgr_agent_node_id }}" --server "{{ vpnmgr_agent_server_addr }}"
stderr_logfile = /var/log/supervisor/%(program_name)s_error.log
stdout_logfile = /var/log/supervisor/%(program_name)s.log
autostart=true
autorestart=true
